# systemd/monitor-test.service
[Unit]
Description=Monitor process "test"
After=network.target

[Service]
Type=simple
User=monitoring
Group=monitoring
ExecStart=/usr/local/bin/monitor_test.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
